# config/systemd/dune-artificial-exchange-bot.service
[Unit]
Description=Dune artificial Exchange service
Documentation=file:///opt/DuneAwakeningSelfHost/docs/artificial-exchange.md
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/opt/DuneAwakeningSelfHost
EnvironmentFile=/opt/DuneAwakeningSelfHost/.env
Environment=PYTHONUNBUFFERED=1
ExecStartPre=/opt/DuneAwakeningSelfHost/scripts/build-exchange-catalog.py
ExecStart=/opt/DuneAwakeningSelfHost/scripts/artificial-exchange-bot.py --loop
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target
