# config/systemd/dune-lan-reflection.service
[Unit]
Description=Dune Awakening LAN reflection host route support
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/DuneAwakeningSelfHost
ExecStart=/opt/DuneAwakeningSelfHost/scripts/setup-lan-reflection.sh

[Install]
WantedBy=multi-user.target
