# config/systemd/dune-replica-snapshot.service
[Unit]
Description=Dune remote Postgres replica snapshot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=dune
Group=dune
WorkingDirectory=/home/keith/Documents/code/DuneAwakeningSelfHost
Environment=DUNE_REPLICA_SNAPSHOT_KEEP_HOURS=48
ExecStart=/home/keith/Documents/code/DuneAwakeningSelfHost/scripts/replica-snapshot.sh /home/keith/Documents/code/DuneAwakeningSelfHost/.env
