# config/systemd/dune-backup-offsite.timer
[Unit]
Description=Run DASH offsite backup sync

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
