# config/systemd/dune-map-watchdog.service
[Unit]
Description=Dune map watchdog
After=docker.service
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/opt/DuneAwakeningSelfHost
Environment=COMPOSE_FILES=compose.yaml:compose.allmaps.yaml
Environment=DUNE_WATCH_INTERVAL=30
Environment=DUNE_WATCH_RECOVERY_WAIT=180
Environment=DUNE_WATCH_COOLDOWN=300
Environment=DUNE_WATCH_STARTUP_GRACE=300
Environment=DUNE_WATCH_LOCK_DIR=/tmp/dune-map-watchdog.lock
ExecStart=/opt/DuneAwakeningSelfHost/scripts/watch-maps.sh /opt/DuneAwakeningSelfHost/.env
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
