# config/systemd/dune-full-farm.service
[Unit]
Description=Dune full 30-map farm startup
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
# Run as the owner of the docker compose CLI plugin
# (~/.docker/cli-plugins/docker-compose). Without User=, systemd runs the
# script as root, and root's docker doesn't find the plugin, so `docker
# compose -f ...` becomes `docker -f ...` and exits 125 with
# "unknown shorthand flag: 'f'". Set DUNE_SERVICE_USER in the installer
# environment to override.
User=DUNE_SERVICE_USER_PLACEHOLDER
Group=DUNE_SERVICE_USER_PLACEHOLDER
WorkingDirectory=/opt/DuneAwakeningSelfHost
ExecStart=/opt/DuneAwakeningSelfHost/scripts/start-full-warm-pool.sh /opt/DuneAwakeningSelfHost/.env
ExecStop=/opt/DuneAwakeningSelfHost/scripts/stop-full-warm-pool.sh /opt/DuneAwakeningSelfHost/.env
RemainAfterExit=yes
TimeoutStartSec=1800
TimeoutStopSec=1800

[Install]
WantedBy=multi-user.target
